# vision-frontyard.service
[Unit]
Description=Vision Frontyard Camera (port 5000)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=KA_JET
WorkingDirectory=/home/KA_JET/robotics/jetson-vision
ExecStart=/home/KA_JET/jetson_yolo_gpu/bin/python3 /home/KA_JET/robotics/jetson-vision/Frontyard_Logging.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vision-frontyard
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
